# etc/systemd/system/ledmgmt.service
[Unit]
Description=LED Management

[Service]
Type=idle
ExecStart=/usr/local/nanogate/ledmgmt

[Install]
WantedBy=multi-user.target
